# supervisord.conf
[supervisord]
nodaemon=true

[program:unicorn]
command=bundle exec unicorn -c /etc/unicorn.rb

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
